# roles/restic/templates/restic-prune.timer.j2
[Unit]
Description=Weekly restic prune

[Timer]
OnCalendar={{ restic_prune_schedule }}
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
